# systemd/howdy-selinux-install.service
[Unit]
Description=Install Howdy SELinux policy
ConditionSecurity=selinux
ConditionPathExists=/usr/libexec/howdy-selinux-setup
After=local-fs.target selinux-autorelabel-mark.service selinux-policy-migrate-local-changes.service
RequiresMountsFor=/etc/selinux /var/lib/selinux

[Service]
Type=oneshot
ExecStart=/usr/libexec/howdy-selinux-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
