# systemd/howdy-boot-reenable.service
[Unit]
Description=Disable Howdy for boot login, re-enable after user authenticates
After=local-fs.target
Before=display-manager.service
DefaultDependencies=no

[Service]
Type=simple
ExecStartPre=/bin/bash -c "\
    [ -f /etc/howdy/config.ini ] || exit 0; \
    if grep -q '^disabled' /etc/howdy/config.ini; then \
        sed -i 's/^disabled.*/disabled = true/' /etc/howdy/config.ini; \
    else \
        sed -i '/^\\[core\\]/a disabled = true' /etc/howdy/config.ini; \
    fi"
ExecStart=/usr/libexec/howdy-reenable
TimeoutStartSec=310

[Install]
WantedBy=multi-user.target
